# supervisord.conf
[program:apache2]
command=apachectl -D FOREGROUND

[program:sshd]
command=/usr/sbin/sshd -D

[supervisord]
nodaemon=true

[supervisorctl]
